# ntfy/rootfs/etc/s6-overlay/s6-rc.d/ntfy/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start the ntfy server
# ==============================================================================

bashio::log.info "Starting ntfy server..."
exec ntfy serve --config /etc/ntfy/server.yml
