# omnivore/rootfs/etc/services.d/redis/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Omnivore
# Runs Redis for Omnivore
# ==============================================================================

# Make sure the redis directory exists
mkdir -p /data/omnivore/redis

# Start Redis with data directory in the persistent storage
bashio::log.info "Starting Redis server..."
redis-server --dir /data/omnivore/redis --appendonly yes
