# vector/rootfs/etc/s6-overlay/s6-rc.d/vector/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start the Vector observability pipeline
# ==============================================================================

declare config_mode
declare config_dir

config_mode=$(bashio::config 'config_mode')
config_dir=$(bashio::config 'config_dir')

if [ "${config_mode}" = "directory" ]; then
    bashio::log.info "Starting Vector with --config-dir ${config_dir}"
    exec vector --config-dir "${config_dir}" --watch-config
else
    VECTOR_CONFIG_PATH="/config/vector/vector.yaml"
    bashio::log.info "Starting Vector with --config ${VECTOR_CONFIG_PATH}"
    exec vector --config "${VECTOR_CONFIG_PATH}" --watch-config
fi
